# etc/crontab
PATH=/bin:/sbin:/usr/bin:/usr/local/bin
* * * * * php /var/www/pimcore/bin/console maintenance -vv >> /var/log/cron.log 2>&1
#
